# Library/LaunchAgents/com.rkn.cleanexfat.plist
<?xml version="1.0" encoding="UTF-8"?>
<plist version="1.0">
  <dict>
    <key>Label</key><string>com.rkn.cleanexfat</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/YOUR_USERNAME/bin/clean-exfat.sh</string>
    </array>
    <key>StartInterval</key><integer>1800</integer>
    <key>RunAtLoad</key><true/>
    <key>StandardOutPath</key><string>/tmp/cleanexfat.out</string>
    <key>StandardErrorPath</key><string>/tmp/cleanexfat.err</string>
  </dict>
</plist>
